# packer/files/puma.service
[Unit]
Description=Puma Service
After=multi-user.target
[Service]
Type=idle
ExecStart=~/reddit/puma -d
[Install]
WantedBy=multi-user.target